# php/.provision/modules/ruby/templates/mailcatcher.conf.erb
start on runlevel [2345]

script
    exec /usr/local/bin/mailcatcher --http-ip=0.0.0.0
end script